# Service_Gateway/Led Status Gateway/Internet_Led_Status/internet-led.service
[Unit]
Description=Internet LED Status Indicator
After=network.target

[Service]
ExecStart=/usr/local/bin/internet-led-status.sh
Restart=always
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target